# ops/digitalocean/systemd/alpha-research-api.service
[Unit]
Description=Alpha Research API
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/alpha-research/repo
Environment=PORT=8787
Environment=DATASET_INSTANCE_ROOT=/srv/alpha-research/data/instances
ExecStart=/usr/bin/node apps/api/dist/server.js
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
