# ops/digitalocean/systemd/alpha-research-frontend.service
[Unit]
Description=Alpha Research Frontend Preview Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/alpha-research/repo/apps/frontend
Environment=VITE_API_BASE_URL=https://api.example.com
ExecStart=/usr/bin/npm exec vite preview -- --host 0.0.0.0 --port 4173
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
